# dist/etc/init.d/cernvm-appliance-agent
#!/bin/sh
# Copyright (c) 2013 CERN
# All rights reserved.
#
# Author: Vasilis Nicolaou, 2013
#
# /etc/init.d/cernvm-appliance-agent
#   and its symbolic link
# /usr/sbin/cernvm-appliance-agent

### BEGIN INIT INFO
# Provides:          cernvm-appliance-agent
# Required-Start:    $network $local_fs $remote_fs 
# Required-Stop:     
# Default-Start:     5
# Default-Stop:      0 1 2 3 4 6
# X-Interactive: true
# Short-Description: Start update sync in boot, act and then die
# Description:       Executes the program that runs the post update actions 
### END INIT INFO

. /etc/init.d/functions 
start() {
        echo -n "Starting cern-vm appliance agent"
        [ -d '/var/log/cernvm-appliance-agent' ] || {
            mkdir -p /var/log/cernvm-appliance-agent
            chown cernvm-appliance-agent /var/log/cernvm-appliance-agent
        }
        [ -f '/var/log/cernvm-appliance-agent/service.log' ] || {
            touch /var/log/cernvm-appliance-agent/service.log
            chown cernvm-appliance-agent /var/log/cernvm-appliance-agent/service.log
        }
	[ -d '/var/lib/cernvm-appliance-agent' ] || {
            mkdir -p /var/lib/cernvm-appliance-agent
            chown cernvm-appliance-agent /var/lib/cernvm-appliance-agent
        }
	[ -f /var/lib/cernvm-appliance-agent/password_set ] && {
	    rm -f /etc/cernvm-appliance-agent/conf.d/*.listen
            ln -s listen.d/all.listen /etc/cernvm-appliance-agent/conf.d/all.listen
        }
        $httpd_BIN -f "$CONFIG_FILE" -k start 2>/var/log/cernvm-appliance-agent/service.log
        if [ $? -eq 0 ]; then
            echo_success
            echo
        else
            echo_failure
            echo
        fi
}

stop() {
        echo -n "Stopping cern-vm appliance agent "
        if [ ! -f /var/run/cernvm-appliance-agent.pid ]; then
            echo "Service is not running"
            echo_failure
            echo
            exit 3
        fi
        
        kill $(cat /var/run/cernvm-appliance-agent.pid)
        if [ $? -eq 0 ]; then
            echo_success
            echo /var/run/cernvm-appliance-agent.pid
            rm /var/run/cernvm-appliance-agent.pid
        else
            echo_failure
            echo
        fi
        exit 0
            
}


# Check for missing binaries
httpd_BIN="/usr/sbin/httpd"
CONFIG_FILE="/etc/cernvm-appliance-agent/conf/httpd-appliance.conf"
test -x $httpd_BIN || { 
        echo "$update_BIN not installed";
        if [ "$1" = "stop" ]; then exit 0;
        else exit 5; fi; }
        
[ -f "$CONFIG_FILE" ] || { 
        echo "$CONFIG_FILE missing";
        if [ "$1" = "stop" ]; then exit 0;
        else exit 5; fi; }

 case "$1" in
    start)
        start
        exit "$?"
        ;;
    stop)
        stop
        # Remember status and be verbose
        exit 3
        ;;
    restart)
        [ ! -f "/var/run/cernvm-appliance-agent.pid" ] && exit 0
        $0 stop
        $0 start
        exit "$?"
        # Remember status and be quiet
        
        ;;
    reload)

        ## Otherwise if it does not support reload:
        exit 3
        ;;
    status)
        echo -n "Checking for service $0: "
        if [ -f "/var/run/cernvm-appliance-agent.pid" ]; then
            echo "Running"
            exit 0
            
        else
            echo "Not running"
            exit 3
        fi
        # 0 - service up and running
        # 1 - service dead, but /var/run/  pid  file exists
        # 2 - service dead, but /var/lock/ lock file exists
        # 3 - service not running (unused)
        # 4 - service status unknown :-(
        # 5--199 reserved (5--99 LSB, 100--149 distro, 150--199 appl.)

        exit 3
        
        ;;
    ?)
        # If no parameters are given, print which are avaiable.
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1
        ;;
esac
